# open-app-filter/files/appfilter.init.bak
#!/bin/sh /etc/rc.common

START=96

start() {
	local enabled
	config_get_bool enabled global enable 0
	if [ "$1" == "1" ] ;then
	{
		insmod oaf
		gen_class.sh /etc/appfilter/feature.cfg
		appfilter.sh
		logger appfilter "enable appfilter"
	}
}

stop() {
	rmmod oaf
	logger appfilter "stop appfilter"	
}

