# open-app-filter/files/appfilter.init
#!/bin/sh /etc/rc.common

START=96
start() {
	insmod oaf
	gen_class.sh /etc/appfilter/feature.cfg
	appfilter.sh
}

stop() {
	rmmod oaf
	echo "stop appfilter"	
}
